# chordial.socket
[Unit]
Description=Chordial

[Socket]
ListenStream=9077

[Install]
WantedBy=chordial.target